# .config/autostart/audiorelay.desktop
[Desktop Entry]
Exec=audiorelay
Name=audiorelay
Type=Application
Version=1.0
