# .config/autostart/emacs.desktop
[Desktop Entry]
Exec=emacs
Name=emacs
Type=Application
Version=1.0
