# .config/autostart/files.desktop
[Desktop Entry]
Exec=nautilus
Name=Nautilus
Type=Application
Version=1.0
