# .config/autostart/io.github.vikdevelop.SaveDesktop.Updater.desktop
[Desktop Entry]
Name=SaveDesktop Native Updater
Type=Application
Exec=python3 /home/kayon/.local/share/savedesktop/native_installer.py --update
