# .config/autostart/keepassxc.desktop
[Desktop Entry]
Exec=keepassxc
Name=keepassxc
Type=Application
Version=1.0
