# .config/autostart/quick.desktop
[Desktop Entry]
Exec=st -e tmux
Name=quick
Type=Application
Version=1.0
