# .config/autostart/slack.desktop
[Desktop Entry]
Exec=slack --gtk-version=3
Name=Slack
Type=Application
Version=1.0
