# .config/autostart/slack.desktop
[Desktop Entry]
Exec=flatpak run com.slack.Slack
Name=Slack
Type=Application
Version=1.0
