# .config/autostart/steam.desktop
[Desktop Entry]
Exec=steam
Name=Steam
Type=Application
Version=1.0
