# .config/autostart/thunar.desktop
[Desktop Entry]
Exec=doublecmd
Name=doublecmd
Type=Application
Version=1.0
