# .config/autostart/virt-manager.desktop
[Desktop Entry]
Exec=virt-manager
Name=virt-manager
Type=Application
Version=1.0
