# .config/autostart/whatsapp.desktop
[Desktop Entry]
Exec=chromium --user-data-dir=/home/kayon/chrome/whats --app=https://web.whatsapp.com
Name=whatsapp
Type=Application
Version=1.0
