# .config/autostart/xrdb.desktop
[Desktop Entry]
Exec=xrdb ~/.Xresources
Name=xrdb
Type=Application
Version=1.0
