# .config/autostart/zen-browser.desktop
[Desktop Entry]
Exec=zen-browser
Name=zen-browser
Type=Application
Version=1.0
