# autostart/.config/autostart/discord.desktop
[Desktop Entry]
Exec=flatpak run com.discordapp.Discord
Name=Discord
Type=Application
Version=1.0
