# autostart/.config/autostart/dolphin.desktop
[Desktop Entry]
Exec=dolphin
Name=Dolphin
Type=Application
Version=1.0
