# autostart/.config/autostart/floorp.desktop
[Desktop Entry]
Exec=librewolf
Name=Librewolf
Type=Application
Version=1.0
