# autostart/.config/autostart/quick.desktop
[Desktop Entry]
Exec=alacritty --option font.size=14 --class quick -e tmux
Name=quick
Type=Application
Version=1.0
