# autostart/.config/autostart/quick.desktop
[Desktop Entry]
Exec=alacritty -e tmux --class quick
Name=quick
Type=Application
Version=1.0
