# autostart/.config/autostart/redshift.desktop
[Desktop Entry]
Exec=redshift
Name=redshift
Type=Application
Version=1.0
