# autostart/.config/autostart/slack.desktop
[Desktop Entry]
Exec=slack
Name=Slack
Type=Application
Version=1.0
