# autostart/.config/autostart/warp-cli.desktop
[Desktop Entry]
Exec=warp-cli disconnect
Name=warp-cli
Type=Application
Version=1.0
